# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:api]
command = uvicorn server:app --host 0.0.0.0 --port 8000
directory = /app
user = root
autostart = true
autorestart = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/fd/2
stderr_logfile_maxbytes = 0
environment = PYTHONUNBUFFERED=1

[program:worker]
command = python worker.py
directory = /app
user = root
autostart = true
autorestart = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/fd/2
stderr_logfile_maxbytes = 0
environment = PYTHONUNBUFFERED=1

[program:telegram-server]
command = python -m telegram_bot.server
directory = /app
user = root
autostart = true
autorestart = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/fd/2
stderr_logfile_maxbytes = 0
environment = PYTHONUNBUFFERED=1